# systemd/powermgr-sync-state.service
[Unit]
Description=Power Manager State Synchronization (Backup Node)
Documentation=https://github.com/your-repo/powermgr
After=network.target

[Service]
Type=oneshot
User=pi
Group=pi
# Sync state file from primary node
# Replace PRIMARY_NODE_IP with actual IP from config
ExecStart=/usr/bin/rsync -a --timeout=30 pi@PRIMARY_NODE_IP:/mnt/ramdisk/powermgr_state.json /mnt/ramdisk/
StandardOutput=journal
StandardError=journal
SyslogIdentifier=powermgr-sync

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/mnt/ramdisk
