# systemd/powermgr-metrics.timer
[Unit]
Description=Power Manager Daily Metrics Timer
Documentation=https://github.com/your-repo/powermgr
Requires=powermgr-metrics.service

[Timer]
# Summer schedule (May-Oct): weekdays at 20:15
# Winter schedule (Nov-Apr): weekdays at 09:15 and 21:15
OnCalendar=May..Oct Mon..Fri 20:15:00
OnCalendar=Nov..Apr Mon..Fri 09:15:00
OnCalendar=Nov..Apr Mon..Fri 21:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
